# packages.list
# Packages to bundle into the ImmortalWrt x86_64 image via ImageBuilder.
# One package per line. '#' starts a comment.
#
# Lines added/removed here trigger a new build (this file is NOT in the
# workflow's paths-ignore list, unlike scripts/*.ps1).
#
# IMPORTANT: these are .ipk names from the ImmortalWrt official + kiddin9
# repositories — NOT source package names. No CONFIG_PACKAGE_ prefix.

# --- core LuCI + Chinese localisation ---
# ImmortalWrt 24.10.2 bakes zh-cn translations into the main luci-app-*
# packages themselves, so separate luci-i18n-*-zh-cn sub-packages either
# don't exist (luci-i18n-opkg-zh-cn) or file-clash on the .lmo translation
# files when installed alongside the main package. We rely on the main
# packages' built-in Chinese strings.
luci
luci-ssl-openssl
luci-i18n-base-zh-cn
luci-i18n-firewall-zh-cn

# --- everyday LuCI apps ---
luci-app-firewall
luci-app-opkg
luci-app-ttyd
luci-app-upnp

# --- Argon theme (modern look + tweak panel) ---
# zh-cn for argon-config is built into luci-app-argon-config, no separate i18n.
luci-theme-argon
luci-app-argon-config

# --- proxy: PassWall + PassWall 2 + cores ---
# Same i18n note: PassWall and PassWall 2 ship zh-cn baked in.
# luci-app-passwall2 comes from kiddin9 mirror (not in upstream feed).
luci-app-passwall
luci-app-passwall2
xray-core
sing-box
hysteria
shadowsocks-rust-sslocal

# --- proxy: OpenClash (Mihomo core fetched at first run) ---
# zh-cn baked into luci-app-openclash.
luci-app-openclash

# --- mesh VPN: Tailscale ---
# LuCI front-end (luci-app-tailscale) is intentionally NOT installed:
# kiddin9's package re-ships /etc/init.d/tailscale and /etc/config/tailscale
# which clash with the upstream tailscale daemon package. Configure via SSH:
#   ssh root@<router>
#   tailscale up      # follow the URL it prints to authenticate
tailscale

# --- mesh VPN: ZeroTier ---
zerotier
luci-app-zerotier

# --- Cloudflare Tunnel / Zero Trust client (CLI only) ---
cloudflared

# --- KMS server for LAN Windows / Office activation ---
vlmcsd
luci-app-vlmcsd

# --- SNMP daemon for Prometheus / Zabbix / LibreNMS monitoring
#     (LuCI front-end omitted — clashes with daemon init script, same as tailscale) ---
snmpd

# --- DNS + DHCP must be the full version (PassWall / OpenClash need ipset support) ---
dnsmasq-full

# --- useful base utilities ---
curl
htop
nano
wget-ssl
bash
ca-bundle
unzip
